# docker-files/supervisord.conf
[supervisord]
nodaemon=true

[program:verifier-api]
command=node /rsk-contract-verifier/dist/contract-verifier-api.js
autostart=true
autorestart=true
